# etc/supervisor/conf.d/check.conf
[program:check]
command=bash -c "sleep 60 && exec /app/venv/bin/weblate check --deploy"
autorestart = false
exitcodes = 0,1
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
