# etc/supervisor/conf.d/check.conf
[program:check]
command = /app/venv/bin/weblate check --deploy
autorestart = false
stdout_events_enabled=true
stderr_events_enabled=true
exitcodes = 0,1
